# connectionmap.service
[Unit]
Description=Connection Map Service

[Service]
ExecStart=/bin/sh /var/www/init.sh

[Install]
WantedBy=multi-user.target
